# root/etc/s6-overlay/s6-rc.d/init-xmind-config/up
/etc/s6-overlay/s6-rc.d/init-xmind-config/run